# deploy/the-full-web.service
[Unit]
Description=the_full_web Next.js service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/the_full_web
Environment=NODE_ENV=production
Environment=HOME=/home/ec2-user
Environment=PATH=/home/ec2-user/.nvm/versions/node/v22.16.0/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ec2-user/.nvm/versions/node/v22.16.0/bin/npm start
Restart=always
RestartSec=5
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
